# ops/lightsail/systemd/pred-infra-pipeline.timer
[Unit]
Description=Run pred-infra pipeline every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=pred-infra-pipeline.service

[Install]
WantedBy=timers.target
